# deploy/systemd/whrkhldsb-next.service.example
[Unit]
Description={{SITE_NAME}} Next.js application
After=network-online.target postgresql.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory={{APP_DIR}}
EnvironmentFile={{ENV_FILE}}
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PATH={{SYSTEMD_PATH}}
ExecStart={{NPM_BIN}} run start
Restart=always
RestartSec=5
User={{APP_USER}}
Group={{APP_USER}}

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/{{APP_SLUG}} /var/backups/{{APP_SLUG}} /tmp {{APP_DIR}}/.next
PrivateTmp=true

# Resource limits (tune for your VPS)
MemoryMax=768M
MemoryHigh=640M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{SERVICE_PREFIX}}-next

[Install]
WantedBy=multi-user.target
