# deploy/systemd/whrkhldsb-ssh-ws.service.example
[Unit]
Description={{SITE_NAME}} SSH WebSocket proxy
After=network-online.target {{SERVICE_PREFIX}}-next.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory={{APP_DIR}}
EnvironmentFile={{RUNTIME_ENV_FILE}}
Environment=NODE_ENV=production
Environment=SSH_WS_PORT={{SSH_WS_PORT}}
Environment=PATH={{SYSTEMD_PATH}}
ExecStart={{NODE_BIN}} {{APP_DIR}}/dist/ssh-ws-proxy.js
Restart=always
RestartSec=5
User={{APP_USER}}
Group={{APP_USER}}

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp
PrivateDevices=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RemoveIPC=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @cpu-emulation @debug @obsolete
SystemCallErrorNumber=EPERM
CapabilityBoundingSet=

# Resource limits
MemoryMax=256M
MemoryHigh=192M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{SERVICE_PREFIX}}-ssh-ws

[Install]
WantedBy=multi-user.target
