# deploy/systemd/whrkhldsb-next.service.example
[Unit]
Description={{SITE_NAME}} Next.js application
After=network-online.target postgresql.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory={{APP_DIR}}
EnvironmentFile={{RUNTIME_ENV_FILE}}
Environment=NODE_ENV=production
Environment=PORT={{NEXT_PORT}}
Environment=PATH={{SYSTEMD_PATH}}
ExecStart={{NODE_BIN}} {{APP_DIR}}/dist/server.js
Restart=always
RestartSec=5
User={{APP_USER}}
Group={{APP_USER}}

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
CapabilityBoundingSet=
PrivateTmp=true

ReadWritePaths=/var/lib/{{APP_SLUG}} /var/backups/{{APP_SLUG}} /tmp {{APP_DIR}}/.next {{APP_DIR}}/uploads {{APP_DIR}}/storage {{APP_DIR}}/downloads {{APP_DIR}}/backups {{APP_DIR}}/logs{{QUICK_SERVICE_READ_WRITE_PATHS}}

# Resource limits (tune for your VPS)
MemoryMax=768M
MemoryHigh=640M
CPUQuota=200%
TasksMax=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{SERVICE_PREFIX}}-next

[Install]
WantedBy=multi-user.target
