# deploy/systemd/whrkhldsb-ssh-ws.service.example
[Unit]
Description={{SITE_NAME}} SSH WebSocket proxy
After=network-online.target {{SERVICE_PREFIX}}-next.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory={{APP_DIR}}
EnvironmentFile={{ENV_FILE}}
Environment=NODE_ENV=production
Environment=PATH={{SYSTEMD_PATH}}
ExecStart={{NPM_BIN}} run start:ssh-ws
Restart=always
RestartSec=5
User={{APP_USER}}
Group={{APP_USER}}

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp
PrivateTmp=true

# Resource limits
MemoryMax=256M
MemoryHigh=192M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{SERVICE_PREFIX}}-ssh-ws

[Install]
WantedBy=multi-user.target
